# bare-metal/zfs/backup/k3s/k3s-volume-backup.timer
[Unit]
Description=Run k3s-volume backup every day

[Timer]
OnCalendar=*-*-* 01:00:00
Persistent=true

[Install]
WantedBy=timers.target