# bare-metal/backups/etcd/etcd-backup.service
[Unit]
Description=etcd backup

[Service]
Type=oneshot
User=ubuntu
ExecStart=/usr/local/bin/etcd-backup.sh